# crates/chmlfrp-tui/assets/chmlfrp-tui.service
[Unit]
Description=ChmlFrp Launcher TUI Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/chmlfrp-tui --daemon
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
